# systemd/ingestion.service
[Unit]
Description=Unified Ingestion Service (Reading + Brain Sync)
Wants=ingestion.timer
After=network.target

[Service]
Type=oneshot
User=server
WorkingDirectory=/home/server/software/observability-hub
EnvironmentFile=/home/server/software/observability-hub/.env
ExecStart=/home/server/software/observability-hub/dist/ingestion
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
