# systemd/mcp-telemetry.service
[Unit]
Description=MCP Telemetry Server
After=network.target

[Service]
Type=simple
User=server2
WorkingDirectory=/home/server2/software/observability-hub
EnvironmentFile=/home/server2/software/observability-hub/.env
ExecStart=/home/server2/software/observability-hub/bin/mcp_telemetry
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
