# systemd/openbao.service
[Unit]
Description=OpenBao Secret Manager
Documentation=https://openbao.org/
After=network.target

[Service]
Type=simple
User=server2
WorkingDirectory=/home/server2/software/observability-hub
# Run bao directly from the system path
ExecStart=/usr/bin/bao server -config=config/openbao/server.hcl
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
# Prevent swapping of memory to disk (crucial for secrets security)
LimitMEMLOCK=infinity

[Install]
WantedBy=multi-user.target
