# systemd/proxy.service
[Unit]
Description=Observability Hub Proxy Server
After=network.target postgresql.service

[Service]
Type=simple
User=server
WorkingDirectory=/home/server/software/observability-hub
# Pointing to the binary we will build
ExecStart=/home/server/software/observability-hub/dist/proxy_server
Restart=always
RestartSec=5
# Load environment variables from the root .env
EnvironmentFile=/home/server/software/observability-hub/.env
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
