# systemd/proxy.service
[Unit]
Description=Observability Hub Proxy Server
After=network.target postgresql.service

[Service]
Type=simple
User=server
WorkingDirectory=/home/server/software/observability-hub
EnvironmentFile=/home/server/software/observability-hub/.env
ExecStart=/home/server/software/observability-hub/bin/proxy_server
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
