# systemd/reading-sync.timer
[Unit]
Description=Twice Daily Sync for Reading Analytics

[Timer]
# Run every 12 hours
OnCalendar=*-*-* 00,12:00:00
# Ensure catch-up if the machine was off
Persistent=true
# Prevent thundering herd (30m window)
RandomizedDelaySec=1800

[Install]
WantedBy=timers.target
