# systemd/reading-sync.service
[Unit]
Description=Trigger Reading Analytics Sync
Wants=reading-sync.timer
After=network.target

[Service]
Type=oneshot
User=server
WorkingDirectory=/home/server/software/observability-hub
EnvironmentFile=/home/server/software/observability-hub/.env
ExecStart=/home/server/software/observability-hub/dist/reading-sync
# Standardize logging for journald
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
