# systemd/tailscale-gate.service
[Unit]
Description=Tailscale Funnel Gatekeeper Loop
After=tailscaled.service

[Service]
ExecStart=/usr/local/bin/tailscale_gate.sh
Restart=always
RestartSec=60
# Ensures the script is prioritized during high RAM usage
OOMScoreAdjust=-500

[Install]
WantedBy=multi-user.target
